# systemd-units/sl-server.service
#MIT License
#
#Copyright (c) 2026 worksoup <https://github.com/worksoup/>
#
#Permission is hereby granted, free of charge, to any person obtaining a copy
#of this software and associated documentation files (the "Software"), to deal
#in the Software without restriction, including without limitation the rights
#to use, copy, modify, merge, publish, distribute, sublicense, and/or sell
#copies of the Software, and to permit persons to whom the Software is
#furnished to do so, subject to the following conditions:
#
#The above copyright notice and this permission notice shall be included in all
#copies or substantial portions of the Software.
#
#THE SOFTWARE IS PROVIDED "AS IS", WITHOUT WARRANTY OF ANY KIND, EXPRESS OR
#IMPLIED, INCLUDING BUT NOT LIMITED TO THE WARRANTIES OF MERCHANTABILITY,
#FITNESS FOR A PARTICULAR PURPOSE AND NONINFRINGEMENT. IN NO EVENT SHALL THE
#AUTHORS OR COPYRIGHT HOLDERS BE LIABLE FOR ANY CLAIM, DAMAGES OR OTHER
#LIABILITY, WHETHER IN AN ACTION OF CONTRACT, TORT OR OTHERWISE, ARISING FROM,
#OUT OF OR IN CONNECTION WITH THE SOFTWARE OR THE USE OR OTHER DEALINGS IN THE
#SOFTWARE.
[Unit]
Description=Simple local server: Tieba sign daemon & tracker merger
After=network.target
Wants=network.target

[Service]
Type=simple
DynamicUser=yes
ExecStartPre=/bin/sh -c 'if [ ! -f /var/lib/sl-server/sl-server.toml ]; then install -m 640 /etc/sl-server/sl-server.toml.example /var/lib/sl-server/sl-server.toml; fi'
ExecStart=/usr/bin/simple_local_server -c /var/lib/sl-server/sl-server.toml

NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/sl-server
StateDirectory=sl-server

StandardOutput=journal
StandardError=journal
SyslogIdentifier=simple-local-server

LimitNOFILE=65536
LimitNPROC=4096

[Install]
WantedBy=multi-user.target
